# roles/borg-client/files/borg_create_archive.service
[Unit]
Description=System backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borg_create_archive.sh
