# roles/nfs/templates/bind.mount.j2
[Unit]
Description=Bind mount {{ item.path }} for NFS sharing

[Mount]
What={{ item.path }}
Where=/srv/nfs/{{ item.dirname }}
Options=bind

[Install]
WantedBy=multi-user.target
