# roles/nfs/templates/nfs.automount.j2
[Unit]
Description=Automount for {{ item.where }}

[Automount]
Where={{ item.where }}

[Install]
WantedBy=multi-user.target
