# roles/mailserver/files/teach-spam.timer
[Unit]
Description=Periodically run teach-spam

[Timer]
OnCalendar=daily
