# roles/mailserver/files/teach-spam.service
[Unit]
Description=Train SpamAssassin's Bayesian filter

[Service]
Type=oneshot
ExecStart=/usr/local/teach-spam
