# roles/scollector/files/supervisor-scollector.conf
[program:scollector]

command = /usr/local/bin/scollector --conf="/etc/scollector.toml"
autorestart = true
startsecs = 5
user = stats
